# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
loglevel=info

; PHP-FPM harus start LEBIH DULU dari Nginx (priority lebih kecil = start lebih awal)
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
priority=1
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Nginx start SETELAH PHP-FPM ready (priority lebih besar = start belakangan)
[program:nginx]
command=nginx -g 'daemon off;'
priority=10
autostart=true
autorestart=true
startretries=5
startsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
